# pi_backend/medidispense.service
[Unit]
Description=MediDispense Pi Backend (Flask API + State Machine)
After=network.target
Wants=network.target

[Service]
Type=simple
User=fens402
WorkingDirectory=/home/fens402/dispenser/pi_backend

# venv içindeki Python ile çalışır (Flask ve tüm pip paketleri burada)
ExecStart=/home/fens402/dispenser/pi_backend/.venv/bin/python3 /home/fens402/dispenser/pi_backend/api_server.py

Environment=HOME=/home/fens402
EnvironmentFile=-/home/fens402/dispenser/pi_backend/.env

# Force-kill if it doesn't stop within 5 seconds
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=5

Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal
SupplementaryGroups=video
TimeoutStartSec=30

[Install]
WantedBy=multi-user.target
